# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:sendmail]
startsecs = 0
autorestart = false
command=/etc/init.d/postfix start

[program:php]
command=php-fpm
